# systemd/iot-web.service
[Unit]
Description=IoT Web Server (タブレット/家族UI)
After=network-online.target iot-matter.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/iot-life-support
ExecStart=/home/pi/iot-life-support/venv/bin/uvicorn src.web.app:app --host 0.0.0.0 --port 8000
ExecStartPost=/bin/bash -c 'sleep 3 && /home/pi/iot-life-support/scripts/notify_url.sh &'
Restart=on-failure
RestartSec=5
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
