# systemd/iot-matter.service
[Unit]
Description=IoT Matter Server (P110M接続用)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/iot-life-support
ExecStart=/home/pi/iot-life-support/venv/bin/matter-server --storage-path /home/pi/iot-life-support/data/matter --port 5580 --bluetooth-adapter 0
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
